# task04/10-keyboard.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{bInterfaceProtocol}=="01" RUN+="/usr/sbin/insmod /root/little-penguin-1/task04/keyboard-task04.ko"
